# deploy/systemd/codex2api.service
[Unit]
Description=Codex2API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deployuser
Group=deployuser
WorkingDirectory=/deploy/codex2api/current
EnvironmentFile=/deploy/codex2api/shared/.env
ExecStart=/deploy/codex2api/current/codex2api
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/deploy/codex2api

[Install]
WantedBy=multi-user.target
