# deploy/systemd/codex2api.service
[Unit]
Description=Codex2API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deployuser
Group=deployuser
WorkingDirectory=/deploy/current
EnvironmentFile=/deploy/shared/.env
ExecStart=/deploy/current/codex2api
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/deploy

[Install]
WantedBy=multi-user.target
