# deploy/systemd/codex2api.service
[Unit]
Description=Codex2API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex2api
Group=codex2api
WorkingDirectory=/opt/codex2api/current
EnvironmentFile=/opt/codex2api/shared/.env
ExecStart=/opt/codex2api/current/codex2api
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/codex2api

[Install]
WantedBy=multi-user.target

